# .config/systemd/user/rclone-bisync.timer
[Unit]
Description=Run rclone-bisync every 30 mins.

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
